# ansible/eth2stats/files/eth2stats.service
[Unit]
Description=ETH2 Stats Service

[Service]
ExecStart=/usr/bin/eth2stats_run
Restart=always

[Install]
WantedBy=multi-user.target
