# nimbus/files/nimbus.service
[Unit]
Description=Nimbus Service

[Service]
ExecStart=/usr/bin/nimbus_run
Restart=always
RuntimeMaxSec=21600

[Install]
WantedBy=multi-user.target
